# confs/extracting_talker_api.conf
[program:extracting_talker_api]
command=/home/th3go2/extracting-talker/scripts/deploy.sh
user=th3go2
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/urban/logs/api.log